# private_dot_config/systemd/user/tpm-update.timer
[Unit]
Description=Run tpm-update twice a week

[Timer]
OnCalendar=Mon,Thu
Persistent=true

[Install]
WantedBy=timers.target
