# private_dot_config/systemd/user/zinit-update.timer
[Unit]
Description=Run zinit-update twice a week

[Timer]
OnCalendar=Mon,Thu
Persistent=true

[Install]
WantedBy=timers.target
